# deploy/systemd/fullpcr.service
[Unit]
Description=fullpcr Streamlit GUI
After=network.target

[Service]
Type=simple
User=fullpcr
Group=fullpcr
WorkingDirectory=/opt/fullpcr
EnvironmentFile=/etc/fullpcr/fullpcr.env
ExecStart=/usr/bin/env python -m fullpcr gui --host ${FULLPCR_HOST} --port ${FULLPCR_PORT} --data-dir ${FULLPCR_DATA_DIR}
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=60
UMask=0027
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
